# deploy/haskellito-backend.supervisor.conf
[program:haskellito-backend]
directory=/opt/haskellito/backend
command=/opt/haskellito/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000
user=haskellito
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/haskellito/backend.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PATH="/opt/haskellito/backend/venv/bin:/usr/bin",GHCI_HOME="/tmp",GHCI_TMPDIR="/tmp",COGNITO_AUTH_ENABLED="true",COGNITO_REGION="us-east-1",COGNITO_USER_POOL_ID="us-east-1_dwg2mt4BF",COGNITO_APP_CLIENT_ID="3akv2gdc854btlb5g6057bf67j",CORS_ALLOW_ORIGINS="https://haskellito.com"
